# rootfs/etc/services.d/wmbus_mqtt_bridge/run
#!/usr/bin/with-contenv bashio
set -e

# twardo uruchamiaj nasz skrypt w bashio środowisku
exec /usr/bin/run.sh
